# remote-backup/scripts/pi/restic-backup.service
[Unit]
Description=Nightly remote backup job

[Service]
Type=oneshot
ExecStart=/home/pi/remote-backup/create-backup.sh
User=root
Group=root

